# docker/flask_demo.ini
[supervisord]
nodaemon=true

[program:server]
command=python3 run.py
directory=/root/flask_demo/src
stdout_logfile=/var/log/stdout.log
stderr_logfile=/var/log/stderr.log
autostart=true
autorestart=true
priority=1